# etc/enginedatalogger.service.in
#test
[Unit]
Description=EngineDataLogger
DefaultDependencies=false
After=timers.target
Wants=timers.target

[Service]
Type=simple
Restart=always
ExecStart = @bindir@/enginedatalogger -c @myconfdir@/enginedatalogger.conf
#WorkingDirectory=@runstatedir@/enginedatalogger

#[Install]
WantedBy=multi-user.target
