# etc/enginedatalogger.service.in
# Replace with actual file
[UNIT]
Description=EngineDataLogger
DefaultDependencies=false

[Service]
Type=simple
Restart=always
ExecStart = @bindir@/enginedatalogger -c @myconfdir@/enginedatalogger.conf
WorkingDirectory=@runstatedir@/enginedatalogger
