# etc/enginedatalogger.path.in
[Unit]
Description = Watch /dev/bus/usb/001/ and activate datalogger when ready

[Path]
DirectoryNotEmpty=/dev/bus/usb/001/

[Install]
WantedBy=multi-user.target


