# etc/shutdowncheck.service.in
[Unit]
Description=ShutdownCheck
DefaultDependencies=false

[Service]
Type=simple
ExecStart=@bindir@/shutdowncheck.py
WorkingDirectory=@runstatedir@

[Install]
WantedBy=local-fs.target

